# roles/ir-control-node/templates/ir-control-node.service
[Unit]
Description=IR control service

[Service]
ExecStart=/usr/local/bin/ir-control-node

[Install]
WantedBy=multi-user.target
