# roles/serveo/templates/serveo.timer
[Unit]
Description=Daily restart of Serveo Tunnel Service

[Timer]
OnCalendar=*-*-* 04:00:00
Unit=serveo.service

[Install]
WantedBy=timers.target
