# roles/anime-fetcher/templates/anime-fetcher.service
[Unit]
Description=Anime fetcher

[Service]
ExecStart=/usr/local/bin/anime-fetcher

[Install]
WantedBy=multi-user.target
